# zulip.conf
[supervisord]
nodaemon=true

[program:zulip]
command=/root/zulip/scripts/start-server
autostart=true
autorestart=true
redirect_stderr=true
